# conf/adcrf.service
[Unit]
Description=ADC fast reader

[Service]
ExecStart=/usr/bin/adcrf

[Install]
WantedBy=multi-user.target
